# ddns
#!/sbin/openrc-run

command="/opt/ddns.sh"
supervisor="supervise-daemon"

depend() {
        need net
}
